# luci-app-eqosplus/root/etc/init.d/eqosplus
#!/bin/sh /etc/rc.common
#
# Copyright 2022-2025 lava <byl0561@gmail.com>
#

START=99
USE_PROCD=1

NAME=eqosplus
LOCKFILE="/tmp/$NAME.lock"

# BusyBox flock lacks -w (timeout); poll with -n instead.
# Usage: _flock_wait <fd> <max_seconds>
_flock_wait() {
	local fd=$1 max=$2 i=0
	while [ "$i" -lt "$max" ]; do
		flock -n "$fd" && return 0
		sleep 1
		i=$((i + 1))
	done
	return 1
}

start_instance() {
	procd_open_instance
	procd_set_param command /usr/bin/eqosplusctrl
	procd_set_param respawn
	procd_set_param stderr 1
	procd_close_instance
}

_eqosplus_stop() {
	eqosplus stop
}

_eqosplus_start() {
	if [ "$(uci -q get $NAME.@eqosplus[0].service_enable)" = "1" ]; then
		if [ "x$(uci -q get turboacc.config.sw_flow)" = 'x1' ] || [ "x$(uci -q get turboacc.config.sfe_flow)" = 'x1' ]; then
			logger -t eqosplus "ERROR: turboacc flow offloading conflicts with eqosplus. QoS rules will be bypassed. Disable sw_flow/sfe_flow in turboacc first."
			return 1
		fi
		eqosplus start || { logger -t eqosplus "eqosplus start failed, not launching daemon"; return 1; }
		start_instance
	fi
}

start_service(){
	# FD 200: service-level lock (prevents concurrent start/stop/reload).
	# Must close explicitly — procd calls these functions in-process,
	# so leaked FDs accumulate across service operations.
	# Signal trap ensures FD is closed even if interrupted mid-operation.
	exec 200>"$LOCKFILE"
	trap 'exec 200>&-' INT TERM HUP
	_flock_wait 200 15 || { exec 200>&-; trap - INT TERM HUP; logger -t eqosplus "start_service: could not acquire lock after 15s, skipping"; return 1; }
	_eqosplus_stop
	_eqosplus_start
	exec 200>&-
	trap - INT TERM HUP
}

service_triggers() {
	procd_add_reload_trigger 'eqosplus'
}

reload_service() {
	# Delegate to start_service so procd's _procd_wrapper tracks instance
	# changes. When service_enable='0', no instance is registered and procd
	# kills the old daemon.  The previous implementation bypassed procd's
	# instance lifecycle, leaving the daemon running after disable.
	start_service "$@"
}

stop_service(){
	exec 200>"$LOCKFILE"
	trap 'exec 200>&-' INT TERM HUP
	_flock_wait 200 15 || { exec 200>&-; trap - INT TERM HUP; logger -t eqosplus "stop_service: could not acquire lock"; return 1; }
	_eqosplus_stop
	exec 200>&-
	trap - INT TERM HUP
}
